# recipes-core/static-ip-config/files/static-ip-config.service
[Unit]
Description=static-ip-config Service

[Service]
ExecStart=/usr/bin/static-ip-config.sh

[Install]
WantedBy=multi-user.target

