# recipes-core/armid/files/armid.service
[Unit]
Description=Arm id Service

[Service]
ExecStart=/usr/bin/armid.sh

[Install]
WantedBy=multi-user.target

