# recipes-core/base-files/files/hostname.service
[Unit]
Description=hostname Service

[Service]
ExecStart=/usr/bin/hostname.sh

[Install]
WantedBy=multi-user.target

