# recipes-core/uds2uart/files/uds2uart.service
[Unit]
Description=uds2uart Service

[Service]
ExecStart=/usr/bin/uds2uart.sh

[Install]
WantedBy=multi-user.target

