# recipes-core/loadfpga/files/loadfpga.service
[Unit]
Description=LoadFpga Service

[Service]
ExecStart=/usr/bin/loadfpga.sh

[Install]
WantedBy=multi-user.target

