# recipes-core/uds2spi/files/uds2spi.service
[Unit]
Description=uds2spi Service

[Service]
ExecStart=/usr/bin/uds2spi.sh

[Install]
WantedBy=multi-user.target

